# package/kochiya/luci-app-emmc-info/root/etc/uci-defaults/luci-app-emmc-info
#!/bin/sh

/etc/init.d/uhttpd restart
chmod 0777 /usr/bin/emmcinfo.sh
rm -rf /tmp/luci-*
exit 0



